# .config/crontab
* * * * * cd /var/www/html /usr/local/bin/php classes/Cron/execCrons.php >> /var/log/cron.log 2>&1
