# .config/crontab
* * * * * php /var/www/html/console runCron >> /var/log/cron.log 2>&1
